# config/com.radon.exit-order-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.radon.exit-order-service</string>
    
    <!-- Description -->
    <key>ServiceDescription</key>
    <string>Convex Scavenger Exit Order Service - Places pending exit orders when IB accepts them</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/joemccann/dev/apps/finance/radon/scripts/exit_order_service.py</string>
    </array>
    
    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>
    
    <!-- Also run at load (startup) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/joemccann/dev/apps/finance/radon</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/joemccann/dev/apps/finance/radon</string>
    </dict>
    
    <!-- Log stdout -->
    <key>StandardOutPath</key>
    <string>/Users/joemccann/dev/apps/finance/radon/logs/exit-order-service.out.log</string>
    
    <!-- Log stderr -->
    <key>StandardErrorPath</key>
    <string>/Users/joemccann/dev/apps/finance/radon/logs/exit-order-service.err.log</string>
    
    <!-- Keep alive settings - only run during work hours -->
    <!-- Note: launchd doesn't have built-in market hours, so we handle this in the script -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>
    
    <!-- Timeout -->
    <key>TimeOut</key>
    <integer>120</integer>
</dict>
</plist>
