# config/com.radon.ibc-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.radon.ibc-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/joemccann/dev/apps/finance/radon/vendor/ibc/scripts/displaybannerandlaunch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joemccann/dev/apps/finance/radon</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>APP</key>
        <string>GATEWAY</string>
        <key>TWS_MAJOR_VRSN</key>
        <string>10.44</string>
        <key>IBC_PATH</key>
        <string>/Users/joemccann/dev/apps/finance/radon/vendor/ibc</string>
        <key>IBC_INI</key>
        <string>/Users/joemccann/ibc/config.ini</string>
        <key>LOG_PATH</key>
        <string>/Users/joemccann/ibc/logs</string>
        <key>TWS_PATH</key>
        <string>/Users/joemccann/Applications</string>
        <key>TWOFA_TIMEOUT_ACTION</key>
        <string>restart</string>
    </dict>

    <!-- Mon=1 through Fri=5 at 00:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- IBC/Gateway manage their own lifecycle via AutoRestartTime -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/joemccann/ibc/logs/ibc-gateway.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/joemccann/ibc/logs/ibc-gateway.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
